# echifa/deploy/systemd/echifa-frontend.service
[Unit]
Description=Echifa Next.js Frontend
After=network.target

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/opt/echifa/current/frontend
EnvironmentFile=/etc/echifa/frontend.env
ExecStart=/usr/bin/env npm start -- --hostname 0.0.0.0 --port 3000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
