# echifa/deploy/systemd/echifa-backend.service
[Unit]
Description=Echifa Go Backend
After=network.target

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/opt/echifa/current/backend
EnvironmentFile=/etc/echifa/backend.env
ExecStart=/opt/echifa/current/backend/bin/echifa
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
